# supervisor/conf/apm-server.ini
[program:apm-server]
command=/usr/local/apm-server/apm-server -e -c /usr/local/apm-server/apm-server.yml
autostart=false
startsecs=10
numprocs=1
stdout_logfile=%(here)s/../log/%(program_name)s.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=1MB
stdout_capture_maxbytes=1MB
stderr_logfile=%(here)s/../log/%(program_name)s.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=1MB
stderr_capture_maxbytes=1MB
