# supervisor/conf/kafka.ini
[program:zookeeper]
command=/usr/local/kafka/bin/zookeeper-server-start.sh config/zookeeper.properties
directory=/usr/local/kafka
autostart=false
startsecs=10
numprocs=1
stdout_logfile=%(here)s/../log/%(program_name)s.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=1MB
stdout_capture_maxbytes=1MB
stderr_logfile=%(here)s/../log/%(program_name)s.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=1MB
stderr_capture_maxbytes=1MB
environment=JAVA_HOME=/usr/local/java11
stopasgroup=true
killasgroup=true

[program:kafka]
command=/usr/local/kafka/bin/kafka-server-start.sh config/server.properties
directory=/usr/local/kafka
autostart=false
startsecs=10
numprocs=1
stdout_logfile=%(here)s/../log/%(program_name)s.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=1MB
stdout_capture_maxbytes=1MB
stderr_logfile=%(here)s/../log/%(program_name)s.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=1MB
stderr_capture_maxbytes=1MB
environment=JAVA_HOME=/usr/local/java11
stopasgroup=true
killasgroup=true

[program:kafka-manager]
command=java --add-opens=java.base/sun.nio.ch=ALL-UNNAMED -jar target/kafdrop-3.28.0-SNAPSHOT.jar ----kafka.brokerConnect=localhost:9092 --management.server.port=8999 --server.port=8999 -Xms64m  -XX:PermSize=64m
directory=/root/code/kafdrop
autostart=false
startsecs=10
numprocs=1
stdout_logfile=%(here)s/../log/%(program_name)s.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=1MB
stdout_capture_maxbytes=1MB
stderr_logfile=%(here)s/../log/%(program_name)s.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=1MB
stderr_capture_maxbytes=1MB
environment=JAVA_HOME=/usr/local/java11
